# deb/etc/init/squash-mount
#!/sbin/runscript
# Copyright 1999-2007 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

depend() {
	use netmount
	use nfsmount
	use openvpn.client
}

start() {
	ebegin "Mounting squash filesystems"
	squash-mount.py --mountsquash --mkdir
	squash-mount.py --mountdd
	eend $? "Failed to mount all squash filesystems"
}

stop() {
        ebegin "Umounting squash filesystems"
        squash-mount.py --umountdd
        squash-mount.py --umountsquash         
        eend $? "Failed to umount all squash filesystems"
}

# vim: ts=4 :
