# deb/etc/init/squash-umount.conf
start on runlevel [06]

task

script
	squash-mount.py --umountdd
	squash-mount.py --umountdd
	squash-mount.py --umountdd
        squash-mount.py --umountsquash         
end script
